# mp4/SD_card/partition1/var/lib/dpkg/info/kerneloops-daemon.conffiles
/etc/kerneloops.conf
/etc/default/kerneloops
/etc/init.d/kerneloops
