# mp4/SD_card/partition1/var/lib/dpkg/info/kerneloops-daemon.list
/.
/etc
/etc/kerneloops.conf
/etc/default
/etc/default/kerneloops
/etc/init.d
/etc/init.d/kerneloops
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/kerneloops.8.gz
/usr/share/doc
/usr/share/doc/kerneloops-daemon
/usr/share/doc/kerneloops-daemon/copyright
/usr/share/doc/kerneloops-daemon/changelog.Debian.gz
/usr/sbin
/usr/sbin/kerneloops
/usr/bin
/usr/bin/kerneloops-submit
