# mp4/SD_card/partition1/var/lib/dpkg/info/oem-config.list
/.
/var
/var/lib
/var/lib/localechooser
/var/lib/oem-config
/usr
/usr/lib
/usr/lib/oem-config
/usr/lib/oem-config/post-install
/usr/lib/oem-config/post-install/ubi-reload-keyboard
/usr/lib/ubiquity
/usr/lib/ubiquity/clock-setup
/usr/lib/ubiquity/debian-installer-utils
/usr/lib/ubiquity/debian-installer-utils/register-module.post-base-installer-oem
/usr/lib/ubiquity/tzsetup
/usr/lib/ubiquity/tzsetup/post-base-installer-oem
/usr/sbin
/usr/sbin/oem-config-firstboot
/usr/sbin/oem-config-prepare
/usr/sbin/oem-config-remove
/usr/sbin/oem-config-wrapper
/usr/share
/usr/share/doc
/usr/share/doc/oem-config
/usr/share/doc/oem-config/copyright
/usr/share/doc/oem-config/changelog.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/oem-config-prepare.8.gz
/usr/share/man/man8/oem-config-firstboot.8.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/oem-config
/etc
/etc/init
/etc/init/oem-config.conf
/etc/init.d
/usr/sbin/oem-config
/usr/share/man/man8/oem-config.8.gz
/etc/init.d/oem-config
