# mp4/SD_card/partition1/var/lib/dpkg/info/whoopsie.list
/.
/etc
/etc/default
/etc/default/whoopsie
/etc/init
/etc/init/whoopsie.conf
/etc/init.d
/usr
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/whoopsie.py
/usr/share/doc
/usr/share/doc/whoopsie
/usr/share/doc/whoopsie/copyright
/usr/share/doc/whoopsie/changelog.gz
/usr/bin
/usr/bin/whoopsie
/etc/init.d/whoopsie
