# mp4/SD_card/partition1/etc/init/xillinux-sound.conf
# xillinux-sound: Bring up audio on Xillinux
#

description	"Bring up sound on Xillinux"

start on filesystem 

task

script
exec /usr/local/bin/xillinux-sound
end script
